# supervisord/goproject.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:myserver]
command=/home/andrea/go/bin/romanServer
autostart=true
autorestart=true
redirect_stderr=true